# switch.service
[Unit]
Description=Switch - XMPP AI Assistant Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/switch
ExecStart=%h/switch/scripts/run.sh
Restart=on-failure
RestartSec=10

# Environment
Environment=HOME=%h
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
